# umad-swiftdialog/payload/Library/LaunchDaemons/com.zillow.mdm_check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.zillowgroup.umad.trigger_nag</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/CPE/bin/mdm_check</string>
	</array>
	<key>KeepAlive</key>
	<dict>
			<key>PathState</key>
			<dict>
				<key>/var/tmp/umad/.mdm_check</key>
				<true/>
			</dict>
	</dict>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>